# profile.d/flink_env
FLINK_HOME=$HOME/opt/flink

PATH=$PATH:$FLINK_HOME/bin

export FLINK_HOME PATH
