# profile.d/java_env
JAVA_HOME=$HOME/opt/java

PATH=$PATH:$JAVA_HOME/bin

export JAVA_HOME PATH
