# profile.d/maven_env
MVN_HOME="$HOME/opt/maven"

PATH=$PATH:$MVN_HOME/bin

export MVN_HOME PATH

