# profile.d/sbt_env
SBT_HOME="$HOME/opt/sbt"

PATH=$PATH:$SBT_HOME/bin

export SBT_HOME PATH

