# profile.d/scala_env
SCALA_HOME="$HOME/opt/scala"
PATH=$PATH:$SCALA_HOME/bin

export SCALA_HOME PATH
